# .config/systemd/user/upgrade.timer
[Unit]
Description=Remind users to update the system

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
